# systemd/user/desktop-session.target
[Unit]

Description=User session stuffs
Requires=dbus.socket
AllowIsolate=true


