# systemd/user/pidgin.service
[Unit]
Description=Pidgin

[Service]
ExecStart=/usr/bin/pidgin
Environment=DISPLAY=:0
Restart=on-abort

[Install]
WantedBy=desktop-session.target
