# systemd/user/pidgin.service
[Unit]
Description=Pidgin

[Service]
ExecStart=/usr/bin/pidgin
Restart=on-abort

[Install]
WantedBy=default.target
